# .local/etc/profile.d/java.sh
export JAVA_HOME=/usr/local/lib/java/jdk-11.0.12+7
export PATH=$JAVA_HOME/bin:$PATH

export MAVEN_HOME=/usr/local/lib/apache/maven/apache-maven-3.8.3
export PATH=$MAVEN_HOME/bin:$PATH

export JDTLS_HOME=$HOME/.local/lib/jdtls
