# .local/etc/profile.d/man.sh
Man() {
    /usr/local/bin/nvim -R -c "Man $1 | silent only"
}
