# .local/etc/profile.d/node.sh
export NPM_HOME=$HOME/.local/share/npm
export PATH=$NPM_HOME/bin:$PATH
