# .local/etc/profile.d/ruby.sh
export GEM_HOME=$HOME/.local/share/gem
export PATH=$GEM_HOME/bin:$PATH
