# .local/etc/profile.d/transmission-remote.sh
alias t-basicstats='transmission-remote -n "transmission:transmission" -st'
alias t-del='t-delete'
alias t-delete='transmission-remote -r -t'
alias t-down='cd /var/lib/transmission-daemon/downloads/'
alias t-fullstats='transmission-remote -n "transmission:transmission" -si'
alias t-list='transmission-remote -n "transmission:transmission" -l'
alias t-move='mv /var/lib/transmission-daemon/downloads/* ~/Downloads/torrents'
alias t-reload='sudo service transmission-daemon reload'
alias t-show='watch -n 1 transmission-remote -n "transmission:transmission" -l'
alias t-start='sudo service transmission-daemon start && sleep 3 && . ~/.local/bin/t-add-queue &'
alias t-status='service transmission-daemon status'
alias t-stop='sudo service transmission-daemon stop'
alias t-web='firefox http://localhost:9091'
